# recipes-gnome/hdmi-detect/files/hdmi-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/share/hdmi-detect/hdmi-hotplug.sh"
