# recipes-gnome/hdmi-detect/files/hdmi-detect.service
[Unit]
Description=Check if HDMI display is connected

[Service]
Type=oneshot
ExecStart=/usr/share/hdmi-detect/hdmi-detect.sh

[Install]
