# recipes-navq/navq-persistent-storage/files/navq-network-settings.path
[Unit]
Description=Network Manager connection monitoring

[Path]
PathChanged=/etc/NetworkManager/system-connections

[Install]
WantedBy=multi-user.target
